# templates/usr/lib/systemd/system/irqbalance.service.d/puppet.conf.epp
[Service]
<% @oneshot -%>
Type=oneshot
RemainAfterExit=yes
<% else -%>
Type=simple
<% end -%>
EnvironmentFile=/etc/sysconfig/irqbalance
ExecStart=/usr/sbin/irqbalance --foreground $IRQBALANCE_ARGS
