# templates/usr/lib/systemd/system/irqbalance.service.d/puppet.conf.epp
<%- | Boolean $oneshot = false | -%>
[Service]
<% if $oneshot { -%>
Type=oneshot
RemainAfterExit=yes
<% } else { -%>
Type=simple
<% } -%>
EnvironmentFile=/etc/sysconfig/irqbalance
ExecStart=/usr/sbin/irqbalance --foreground $IRQBALANCE_ARGS
