# templates/etc/systemd/system/irqbalance.service.d/puppet.conf.epp
<%- |
  Boolean $oneshot = false,
  Stdlib::Absolutepath $irqbalance_binary = '/usr/sbin/irqbalance',
  Stdlib::Absolutepath $irqbalance_env_file = '/etc/sysconfig/irqbalance',
| -%>
[Service]
<% if $oneshot { -%>
Type=oneshot
RemainAfterExit=yes
<% } else { -%>
Type=simple
<% } -%>
EnvironmentFile=<%= $irqbalance_env_file %>
ExecStart=
ExecStart=<%= $irqbalance_binary %> --foreground $IRQBALANCE_ARGS
